# external/common/detectAI/detect.service
[Unit]
Description=Detection_gpio.sh
[Service]
Type=forking
ExecStart=/bin/bash -c "/home/orangepi/detectAI/detection_gpio.sh"
Environment=XAUTHORITY=/home/orangepi/.Xauthority
Environment=DISPLAY=:0.0
Restart=always
TimeoutStartSec=0
[Install]
WantedBy=multi-user.target

